# conf/blacklist/update-blacklist.cron
#!/bin/sh
/usr/local/bin/update-blacklist.sh -d
